# ansible/roles/gateway/templates/nftables.service.j2
[Service]
ExecStart=
ExecStart=/usr/sbin/nft -f /etc/nftables.conf
ExecReload=
ExecReload=/usr/sbin/nft -f /etc/nftables.conf
